# infra/ansible/roles/agent/templates/tridot-agent.service.j2
[Unit]
Description=TriDot VPN Node Agent
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/tridot-agent
ExecStart=/opt/tridot-agent/agent --config /etc/tridot/agent.yaml
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
